# etc/api/supervisord.conf
[supervisord]
nodaemon=true

[program:php_server]
command=symfony server:start
